# axiom-flame/systemd/axiom-flame-production.service
[Unit]
Description=AXIOM FLAME API
Documentation=https://github.com/JermaineMerritt-ai/Super-Codex-AI
After=network.target
Wants=network.target

[Service]
Type=simple
User=axiom-flame
Group=axiom-flame
WorkingDirectory=/srv/axiom-flame/packages/api
ExecStart=/srv/axiom-flame/.venv/bin/python app.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10
TimeoutStopSec=30

# Environment
Environment=PORT=8080
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/srv/axiom-flame
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/axiom-flame/storage
ReadWritePaths=/srv/axiom-flame/artifacts
ReadOnlyPaths=/srv/axiom-flame

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=axiom-flame

[Install]
WantedBy=multi-user.target