# axiom-flame/systemd/axiom-flame.service
[Unit]
Description=AXIOM FLAME API
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/axiom-flame/packages/api
ExecStart=/srv/axiom-flame/.venv/bin/python app.py
Restart=always
Environment=PORT=8080

[Install]
WantedBy=multi-user.target