# scripts/codex-surveillance.service
[Unit]
Description=Codex Surveillance & Compliance Capsule
After=network.target codex-execution.service

[Service]
Type=simple
WorkingDirectory=/opt/codex/capsule
ExecStart=/usr/bin/python3 services/surveillance.py \
  --config config/capsule_policy_manifest.yaml
Restart=always
RestartSec=5
User=codex

[Install]
WantedBy=multi-user.target
